# services/climate-monitor-sync.timer
[Unit]
Description=Climate Monitor DB Sync Timer

[Timer]
OnUnitActiveSec=60
Persistent=true

[Install]
WantedBy=timers.target
